# php-fpm/supervisor/reverb.conf
[program:reverb]
environment=APP_ENV="docker",APP_DEBUG="false",APP_PATH="/var/www/lunchbox"
process_name = %(program_name)s_%(process_num)s
command = /usr/local/bin/php /var/www/lunchbox/artisan reverb:start
user = %(ENV_USER)s
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
minfds = 10000
